# infra/systemd/neflo-scan.timer
[Unit]
Description=Run the Neflo deposit scan every 30 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5
Unit=neflo-scan.service

[Install]
WantedBy=timers.target
